# etc/profile.d/apacheds.sh
export PATH=/opt/apacheds/bin:/usr/java/default/bin:$PATH
export JAVA_HOME=/usr/java/default
